# as11loader/10-hdus.rules
SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="1738", ATTRS{idProduct}=="5211", RUN+="/bin/sh -c 'cd /usr/local/lib/hdus;./as11loader'"
SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="3275", ATTRS{idProduct}=="6051", MODE="0664", GROUP="video"
SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="3275", ATTRS{idProduct}=="6081", MODE="0664", GROUP="video"
